# udiskie.service
[Unit]
Description=Handle automounting of usb devices


[Service]
Type=simple
ExecStart=/usr/bin/udiskie --smart-tray


[Install]
WantedBy=default.target
